# daemon/com.claude.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.claude.orchestrator</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/chriscrabtree/.bun/bin/bun</string>
		<string>src/index.ts</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/chriscrabtree/claude-orchestrator</string>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>ThrottleInterval</key>
	<integer>10</integer>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/chriscrabtree/.bun/bin:/usr/local/opt/node@22/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/chriscrabtree</string>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/chriscrabtree/Library/Logs/claude-orchestrator.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/chriscrabtree/Library/Logs/claude-orchestrator.error.log</string>
</dict>
</plist>
